# pi-dns-warden/ops/systemd/pihole-tor-prometheus-watchdog.service.template
[Unit]
Description=Pi-hole Tor Prometheus external watchdog
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@ROOT_DIR@
ExecStart=@ROOT_DIR@/ops/scripts/18-prometheus-watchdog.sh

