# pi-dns-warden/ops/systemd/pihole-tor-prometheus-watchdog.timer.template
[Unit]
Description=Run Pi-hole Tor Prometheus external watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Unit=pihole-tor-prometheus-watchdog.service

[Install]
WantedBy=timers.target
