# pi-dns-warden/ops/systemd/pihole-tor.service.template
[Unit]
Description=Torhole Advanced DNS and operations stack (Docker Compose)
Requires=docker.service pihole-tor-vlans.service
After=docker.service pihole-tor-vlans.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@ROOT_DIR@
RemainAfterExit=yes
ExecStart=/usr/bin/bash @ROOT_DIR@/ops/scripts/22-stack-service.sh up
ExecStop=/usr/bin/bash @ROOT_DIR@/ops/scripts/22-stack-service.sh down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
